# etc/udev/rules.d/95-monitor-switch.rules
ACTION=="change", KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/home/pab/.monitor.sh"
